# addons/node-autoscaler/autoscaler.service
[Unit]
Description=Kubernetes Node Autoscaler for Proxmox
Documentation=https://github.com/jlambert229/k8s-deploy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=owner
Group=owner
WorkingDirectory=/home/owner/Repos/k8s-deploy/addons/node-autoscaler

# Load environment configuration
EnvironmentFile=/home/owner/Repos/k8s-deploy/addons/node-autoscaler/autoscaler.env

# Set PATH to include common binary locations and Homebrew
Environment="PATH=/home/linuxbrew/.linuxbrew/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# Use the virtual environment Python
ExecStart=/home/owner/Repos/k8s-deploy/addons/node-autoscaler/venv/bin/python autoscaler-webhook.py

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=k8s-node-autoscaler

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
